# salt/files/etc/cron.d/highstate
1 * * * * root salt-call --local state.highstate >/dev/null
