# salt/files/etc/cron.d/letsencrypt
1 * * * 0 root /root/letsencrypt.sh >/dev/null
